# monitoring/sensu/sensu-client-v2/supervisor.conf
[supervisord]
nodaemon=true

[program:sensu-client]
command=/etc/init.d/sensu-client start